# file/etc/init.d/pfc

#!/bin/sh /etc/rc.common

START=50

reload()
{
    start
}

restart()
{
    reload
}

start()
{
    /usr/sbin/port_fc start
}

stop()
{
    /usr/sbin/port_fc stop
}

boot()
{
    export SQM_VERBOSITY_MIN=5 # Silence errors
    /usr/sbin/port_fc start
}

